# supervisord.conf
[program:telegram-bot-api]
; Local Telegram Bot API server — removes 20MB file size limit.
; Reads TELEGRAM_API_ID and TELEGRAM_API_HASH from environment automatically.
command=telegram-bot-api --local --dir=/var/lib/telegram-bot-api --temp-dir=/tmp/telegram-bot-api --http-port=8081 --username=telegram-bot-api --groupname=telegram-bot-api
autorestart=true
priority=10
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bot]
; Wait until port 8081 is actually accepting connections before starting the bot.
; This survives slow telegram-bot-api startup on Railway without a fixed sleep.
command=/bin/sh -c 'until python3 -c "import socket; s=socket.socket(); s.settimeout(1); s.connect((\"localhost\",8081)); s.close()" 2>/dev/null; do sleep 2; done && python -m bot'
directory=/app
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
