# supervisord.conf
[program:telegram-bot-api]
; Local Telegram Bot API server — removes 20MB file size limit.
; Reads TELEGRAM_API_ID and TELEGRAM_API_HASH from environment automatically.
command=telegram-bot-api
    --local
    --dir=/var/lib/telegram-bot-api
    --temp-dir=/tmp/telegram-bot-api
    --http-port=8081
autorestart=true
priority=10
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bot]
; Wait for the API server to be ready before starting
command=/bin/sh -c 'sleep 6 && python -m bot'
directory=/app
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
